# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

[program:weaviate]
command=/scripts/start-weaviate.sh
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/var/log/supervisor/weaviate.log
stderr_logfile=/var/log/supervisor/weaviate.log
stdout_logfile_maxbytes=10MB
priority=1

[program:api]
command=/scripts/start-api.sh
autostart=true
autorestart=true
startretries=5
startsecs=15
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api.log
stdout_logfile_maxbytes=10MB
priority=2
